# files/etc/profile.d/scl-enable-rh-git227.sh
#!/bin/bash
source scl_source enable rh-git227
source /opt/rh/rh-git227/root/usr/share/doc/rh-git227-git-2.27.0/contrib/completion/git-prompt.sh
source /opt/rh/rh-git227/root/usr/share/doc/rh-git227-git-2.27.0/contrib/completion/git-completion.bash
