# files/etc/profile.d/scl-enable-rh-git218.sh
#!/bin/bash
source scl_source enable rh-git218
source /opt/rh/rh-git218/root/usr/share/doc/rh-git218-git-2.18.4/contrib/completion/git-completion.bash
